# install_files/udev/99-razer.rules
ACTION!="add", GOTO="razer_end"
SUBSYSTEMS=="usb|input|hid", ATTRS{idVendor}=="1532", GOTO="razer_vendor"
GOTO="razer_end"

LABEL="razer_vendor"

# Mice
ATTRS{idProduct}=="0013|0016|0020|0024|0025|002e|002f|0032|0034|0036|0037|0038|0039|0040|0041|0042|0043|0044|0045|0046|0048|004c|0050|0053|0054|0059|005a|005b|005c|005e|0060|0064|0067|006a|006b|006c|006e|006f|0070|0071|0072|0073|0078|007a|007b|0084|0086|0088", \
    ATTRS{idVendor}=="1532", \
    ENV{ID_RAZER_CHROMA}="1", ENV{RAZER_DRIVER}="razermouse"

# Keyboards
ATTRS{idProduct}=="010d|010e|010f|0111|0113|011a|011b|0201|0202|0203|0204|0205|0207|0208|0209|020f|0210|0211|0214|0216|0217|021a|021e|021f|0220|0221|0224|0225|0226|0227|0228|022a|022b|022d|022f|0232|0233|0234|0235|0237|0239|023a|023b|0240|0241|0243|0245|0246|024a|024c|024d", \
    ATTRS{idVendor}=="1532", \
    ENV{ID_RAZER_CHROMA}="1", ENV{RAZER_DRIVER}="razerkbd"

# Mouse Mats
ATTRS{idProduct}=="0068|0c00|0c01|0c02", \
    ATTRS{idVendor}=="1532", \
    ENV{ID_RAZER_CHROMA}="1", ENV{RAZER_DRIVER}="razermousemat"

# Kraken
ATTRS{idProduct}=="0501|0504|0506|0510|0527", \
    ATTRS{idVendor}=="1532", \
    ENV{ID_RAZER_CHROMA}="1", ENV{RAZER_DRIVER}="razerkraken"

# Accessories - Nommo Chroma, Nommo Pro, Mug, Base Station, HDK
ATTRS{idProduct}=="007e|0517|0518|0f07|0f08|0f09", \
    ATTRS{idVendor}=="1532", \
    ENV{ID_RAZER_CHROMA}="1", ENV{RAZER_DRIVER}="razeraccessory"

# Core
ATTRS{idProduct}=="0215", \
    ATTRS{idVendor}=="1532", \
    ENV{ID_RAZER_CHROMA}="1", ENV{RAZER_DRIVER}="razercore"

# Get out if no match
ENV{ID_RAZER_CHROMA}!="1", GOTO="razer_end"

# Set permissions if this is an input node
SUBSYSTEM=="input|hid", GROUP:="plugdev"

# We're done unless it's the hid node
SUBSYSTEM!="hid", GOTO="razer_end"

# Rebind if needed
SUBSYSTEM=="hid", RUN+="razer_mount $env{RAZER_DRIVER} $kernel"

LABEL="razer_end"
