# install_files/udev/99-razer.rules
ACTION!="add", GOTO="razer_end"
SUBSYSTEMS=="usb|input|hid", ATTRS{idVendor}=="1532", GOTO="razer_vendor"
GOTO="razer_end"

LABEL="razer_vendor"

# Mice
ATTRS{idProduct}=="0013|0016|0020|0024|0025|002e|002f|0032|0034|0036|0037|0038|0039|003e|003f|0040|0041|0042|0043|0044|0045|0046|0048|004c|004f|0050|0053|0054|0059|005a|005b|005c|005e|0060|0062|0064|0065|0067|006a|006b|006c|006e|006f|0070|0071|0072|0073|0077|0078|007a|007b|007c|007d|0080|0083|0084|0085|0086|0088|008a|008c|008d|008f|0090|0091|0094|0095|0098|0099", \
    ATTRS{idVendor}=="1532", \
    ENV{ID_RAZER_CHROMA}="1", ENV{RAZER_DRIVER}="razermouse"

# Keyboards
ATTRS{idProduct}=="010d|010e|010f|0111|0113|011a|011b|011c|0201|0202|0203|0204|0205|0207|0208|0209|020f|0210|0211|0214|0216|0217|021a|021e|021f|0220|0221|0224|0225|0226|0227|0228|022a|022b|022c|022d|022f|0232|0233|0234|0235|0237|0239|023a|023b|023f|0240|0241|0243|0245|0246|024a|024c|024d|024e|0252|0253|0255|0256|0257|0258|0259|025a|025d|025e|0266|0269|026a|026c|026d|026f|0270|0271|0276|0279|028a|028b|0a24", \
    ATTRS{idVendor}=="1532", \
    ENV{ID_RAZER_CHROMA}="1", ENV{RAZER_DRIVER}="razerkbd"

# Kraken
ATTRS{idProduct}=="0501|0504|0506|0510|0527", \
    ATTRS{idVendor}=="1532", \
    ENV{ID_RAZER_CHROMA}="1", ENV{RAZER_DRIVER}="razerkraken"

# Accessories (Speakers, Mousemats, Razer Core, etc)
ATTRS{idProduct}=="0068|007e|0215|0517|0518|0c00|0c01|0c02|0c04|0f07|0f08|0f09|0f12|0f19|0f1d|0f20|0f21|0f26", \
    ATTRS{idVendor}=="1532", \
    ENV{ID_RAZER_CHROMA}="1", ENV{RAZER_DRIVER}="razeraccessory"

# Get out if no match
ENV{ID_RAZER_CHROMA}!="1", GOTO="razer_end"

# Set permissions if this is an input node
SUBSYSTEM=="usb|input|hid", GROUP:="plugdev"

# We're done unless it's the hid node
SUBSYSTEM!="hid|usb", GOTO="razer_end"

# Rebind if needed
SUBSYSTEM=="hid|usb", RUN+="razer_mount $env{RAZER_DRIVER} $kernel"

LABEL="razer_end"
