# daemon/resources/openrazer-daemon.systemd.service.in
[Unit]
Description=Daemon to manage razer devices in userspace
Documentation=man:openrazer-daemon(8)

[Service]
Type=dbus
BusName=org.razer
ExecStartPre=/usr/bin/udevadm settle --timeout=10
ExecStart=##PREFIX##/bin/openrazer-daemon -F

[Install]
WantedBy=default.target
